# tpl/etc/periodic/monthly/maybe-renew
#!/bin/sh
{{if getenv "DOMAIN"}}{{if getenv "LETSENCRYPT" -}}
supervisorctl start certbot
{{- end}}{{end}}

