# tpl/etc/periodic/monthly/maybe-renew
#!/bin/sh
{{if .Env.DOMAIN}}{{if .Env.LETSENCRYPT -}}
supervisorctl start certbot
{{- end}}{{end}}

