# tpl/etc/supervisor.d/certbot.ini
{{if .Env.DOMAIN}}{{if .Env.LETS_ENCRYPT}}
[program:certbot]
command=/usr/bin/maybe-encrypt
autostart=true
autorestart=false
priority=999 ; only run after everything else is up
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
{{end}}{{end}}

