# tpl/etc/supervisor.d/cron.ini
{{if getenv "USE_CRON" | bool -}}
[program:cron]
command=/usr/sbin/crond -f -d 8
autostart=true
autorestart=true
priority=15 ; only run after php & nginx are up
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
{{end}}

