# tpl/etc/supervisor.d/modd.ini
{{if default .Env.MODD_CONF "" -}}
[program:modd]
command=/usr/bin/modd -f {{ .Env.MODD_CONF }} {{ default .Env.MODD_OPTIONS "" }}
directory={{ default .Env.MODD_DIR .Env.CODE_BASE }}
autostart=true
autorestart=true
priority=15 ; only run after php & nginx are up
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
{{end}}

