# tpl/etc/supervisor.d/webhook.ini
{{if default .Env.WEBHOOK_CONF "" -}}
[program:webhook]
command=/usr/bin/webhook -hooks {{ .Env.WEBHOOK_CONF }} -ip 127.0.0.1 -port {{ .Env.WEBHOOK_PORT }} -urlprefix {{ .Env.WEBHOOK_PATH }} {{ .Env.WEBHOOK_OPTS }}
directory={{ default .Env.WEBHOOK_DIR .Env.CODE_BASE }}
user={{ .Env.WEBHOOK_USER }}
autostart=true
autorestart=true
priority=10 ; start with nginx
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
{{end}}

